# kopia/backup-cronjobs/host-setup/99-usb-backup-mount.rules
# Udev rule for auto-mounting USB backup drive
# Place this file in /etc/udev/rules.d/99-usb-backup-mount.rules

# Trigger mount script when USB drive with specific UUID is added or removed
ACTION=="add|remove", ENV{ID_FS_UUID}=="4baeb9e1-2fbd-4985-9ea7-93fbf9f25589", RUN+="/usr/local/bin/usb-backup-mount.sh"
